# root/etc/uci-defaults/40_luci-mihomo
#!/bin/sh
chmod -R +x /etc/init.d/mihomo  /etc/mihomo/* /usr/libexec/mihomo-call

exit 0
